# roles/ejabberd/files/ejabberd-compress-muc-logs.timer
[Unit]
Description=Run daily compression of mod_muc logs for ejabberd

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
FixedRandomDelay=true

[Install]
WantedBy=timers.target
